# payroll-engine/crontab
# payroll-engine/crontab
# Run at 12:00 AM on the 15th and 30th of every month
0 0 15,30 * * python /app/payroll_engine.py --cutoff $(date +\%Y-\%m-%d)
